# mailer/mailer.cron
0 * * * * root /bin/sh /usr/share/mailer/mailer.sh
